# release/etc/systemd/system/veloxd.service.in
[Unit]
Description=Velox userspace TCP/IP daemon
Documentation=https://github.com/SEGV1/velox
Requires=network.target
After=network.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/veloxd --daemon
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=default.target
